# crontab
# 0 7 * * * /app/scripts/indexer
* * * * * /app/scripts/indexer
